# dev_setup/cookbooks/mongodb/templates/default/mongodb.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN
# mongodb.upstart.conf

<%
  @config = node[:mongodb][@server_type]
  @binary = @server_type == :mongos ? 'mongos' : 'mongod'
%>

pre-start script
  <% if @config[:dbpath] %>
  mkdir -p <%= @config[:dbpath] %>
  chown -R mongodb:mongodb <%= @config[:dbpath] %>
  <% end %>
  mkdir -p <%= File.dirname(@config[:logfile]) %>
end script

start on runlevel [2345]
stop on runlevel [06]

script
  ENABLE_MONGODB="yes"

  DAEMON=<%= @binary %>
  NAME=<%= @server_type %>

  # Default defaults. Can be overridden by the /etc/default/$NAME
  CONF=<%= @config[:config] %>
  PIDFILE=<%= @config[:pidfile] %>

  <% case server_type; when 'config_server' %>
  CONFIGSVR=--configsvr
  <% when 'mongos' %>
  CONFIGDB="--configdb <%= @configdb_server_list %>"
  <% else %>
  REPLSET=<%= node[:mongodb][:replica_set] %>
  SHARDSVR=<%= node[:mongodb][:shard_server] %>
  <% end %>

  <% unless server_type == 'mongos' %>
  REST=<%= node[:mongodb][:rest] %>
  SYNCDELAY="--syncdelay <%= node[:mongodb][:syncdelay] %>"
  <% end %>

  if [ -f /etc/default/$NAME ]; then . /etc/default/$NAME; fi

  # Check for daemon bin
  if test ! -x $DAEMON; then
    echo "Could not find $DAEMON"
    exit 0
  fi

  # This will enable the REST Interface for MongoDB if told so
  if test "x$REST" = "xtrue"; then
    REST=--rest
  else
    REST=
  fi

  # This will make the server part of a replica set
  if test "x$REPLSET" = "x"; then
    REPLSET=
  else
    REPLSET="--replSet $REPLSET"
  fi

  # This will make the server participate in sharding
  if test "x$SHARDSVR" = "xtrue"; then
    SHARDSVR=--shardsvr
  else
    SHARDSVR=
  fi

  if [ "x$ENABLE_MONGODB" = "xyes" ]; then exec start-stop-daemon --start --quiet --chuid mongodb --exec  /usr/bin/$DAEMON -- $REST $REPLSET $SHARDSVR $CONFIGSVR $SYNCDELAY $CONFIGDB --config $CONF; fi
end script
